# scripts/launchd/xyz.tanchao.substack-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job: daily Substack notes sync.

  Install to ~/Library/LaunchAgents and load with `launchctl bootstrap`.
  See ../sync-substack-local.sh for setup notes.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>xyz.tanchao.substack-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ctan/.cache/tanchao-substack-sync/scripts/sync-substack-local.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ctan/.cache/tanchao-substack-sync</string>

    <!-- 06:00 local time, daily. Launchd catches up if the Mac was asleep. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ctan/Library/Logs/substack-sync.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ctan/Library/Logs/substack-sync.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
